# systemd/chromecast_controller.service
[Unit]
Description=Chromecast Controller Service

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m castcontroller

