# systemd/chromecast_controller.service
[Unit]
Description=Chromecast Controller Service

[Service]
Type=simple
ExecStart=/opt/Chromecast-Controller/controller.py

